# launchd/ch.invisible.supertag-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ch.invisible.supertag-daily</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Path replaced by install-launchd.sh -->
        <string>/usr/local/bin/supertag</string>
        <string>sync</string>
        <string>index</string>
    </array>

    <!-- Run every 6 hours: midnight, 6 AM, noon, 6 PM -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <!-- Path replaced by install-launchd.sh -->
    <string>/Users/YOUR_USERNAME</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <!-- Path replaced by install-launchd.sh -->
        <string>/Users/YOUR_USERNAME</string>
    </dict>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <!-- Path replaced by install-launchd.sh -->
    <string>/Users/YOUR_USERNAME/.local/state/supertag/logs/supertag-daily.log</string>

    <key>StandardErrorPath</key>
    <!-- Path replaced by install-launchd.sh -->
    <string>/Users/YOUR_USERNAME/.local/state/supertag/logs/supertag-daily.error.log</string>

    <!-- Don't run at load, only on schedule -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Timeout: 30 minutes (sync can take time for large workspaces) -->
    <key>TimeOut</key>
    <integer>1800</integer>

    <!-- Throttle interval: don't restart more than once per hour -->
    <key>ThrottleInterval</key>
    <integer>3600</integer>
</dict>
</plist>
