# launchd/ch.invisible.tana-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ch.invisible.tana-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/fischer/work/DA/KAI/skills/tana/tana-daily</string>
    </array>

    <!-- Run three times daily at 6:00, 12:00, and 18:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- Environment variables for bun/PATH -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/fischer/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/fischer</string>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/fischer/work/DA/KAI/skills/tana</string>

    <!-- Log output (XDG_STATE_HOME/supertag/logs) -->
    <key>StandardOutPath</key>
    <string>/Users/fischer/.local/state/supertag/logs/tana-daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/fischer/.local/state/supertag/logs/tana-daily.error.log</string>

    <!-- Don't run at load, only on schedule -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Timeout: 30 minutes (export can take 15+ min for large workspaces) -->
    <key>TimeOut</key>
    <integer>1800</integer>

    <!-- Restart if the job fails (max 3 retries) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Throttle interval: don't restart more than once per hour -->
    <key>ThrottleInterval</key>
    <integer>3600</integer>
</dict>
</plist>
