# launchd/ch.invisible.tana-webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Identification -->
    <key>Label</key>
    <string>com.pai.tana-webhook</string>

    <!-- Program to Execute -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/fischer/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/fischer/work/DA/KAI/skills/tana/src/cli/tana-webhook.ts</string>
        <string>start</string>
        <string>--port</string>
        <string>3100</string>
        <string>--host</string>
        <string>localhost</string>
        <string>--db-path</string>
        <string>/Users/fischer/work/DA/KAI/skills/tana/test-production.db</string>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/fischer/work/DA/KAI/skills/tana</string>

    <!-- When to Run -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <!-- Restart if crashed -->
        <key>SuccessfulExit</key>
        <false/>
        <!-- Check every 60 seconds if not running -->
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Logging (XDG_STATE_HOME/tana/logs) -->
    <key>StandardOutPath</key>
    <string>/Users/fischer/.local/state/tana/logs/tana-webhook.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/fischer/.local/state/tana/logs/tana-webhook.error.log</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/fischer/.bun/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/fischer</string>
    </dict>

    <!-- Resource Limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>

    <!-- Process Configuration -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <key>Nice</key>
    <integer>0</integer>

    <!-- Throttle restarts if crashing repeatedly -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
